# deb/install/postinst
#!/bin/bash

case "${1}" in
  configure)
    # enable the service to run on on boot
    systemctl enable waggle-nodeid.service
  ;;
esac
